# dab/supervisor/mux.conf
[program:*Multiplex]
command=odr-dabmux conf.mux
directory=/home/pi/dab/
autostart=false
autorestart=false
stderr_logfile=/var/log/supervisor/mux.log
stdout_logfile=/var/log/supervisor/mux.log
